# ops/systemd/theo-mouth-session-bridge.timer
# Theo Mouth session bridge user timer v0.1.0 - Noted by Theo - 2026-06-14.

[Unit]
Description=Run Theo Glass Mouth session bridge every 30 seconds

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
RandomizedDelaySec=5s
AccuracySec=5s
Persistent=false

[Install]
WantedBy=timers.target
